# docent-check-keys.service
[Unit]
Description=Docent VPN Flask backend
After=network.target postgresql.service

[Service]
User=admin
Group=admin
WorkingDirectory=/home/admin/docent-vpn/backend
Environment=
ExecStart=/home/admin/docent-vpn/backend/.venv/bin/python3 check_keys.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
